# project_ocp/ocp_files/dnsmasq.service.j2
[Unit]
Description=DNS caching server.
After=network.target

[Service]
ExecStart=/usr/sbin/dnsmasq -k -p0

[Install]
WantedBy=multi-user.target
